# root/etc/s6-overlay/s6-rc.d/svc-wireguard/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ ${USE_WIREGUARD,,} == "false" ]]; then
    exit 0
else
    wg-quick down wg0
fi
